# examples/emby-proxy.service
[Unit]
Description=Emby reverse proxy helper
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=emby-proxy
Group=emby-proxy
Environment=LISTEN_ADDR=127.0.0.1:8080
Environment=BLOCK_PRIVATE_TARGETS=true
Environment=GODEBUG=http2client=0
ExecStart=/usr/local/bin/emby-proxy
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
